# systemd/pineapplexpress-pcap-recorder.service
[Unit]
Description=PineAppleXpress Rolling PCAPNG Recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pineapple
SupplementaryGroups=wireshark
EnvironmentFile=/etc/default/pineapplexpress-packets
ExecStart=/usr/local/sbin/pineapplexpress-pcap-recorder
Restart=on-failure
RestartSec=3
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
